# deployment/supervisor/zender-worker.conf
; ============================================
; ZENDER QUEUE WORKER SUPERVISOR CONFIGURATION
; ============================================
; This file configures the Laravel queue workers for production.
; 
; Installation:
; 1. Install supervisor: sudo apt install supervisor
; 2. Copy this file to: /etc/supervisor/conf.d/zender-worker.conf
; 3. Update the paths below to match your installation
; 4. Run: sudo supervisorctl reread && sudo supervisorctl update
; 5. Start workers: sudo supervisorctl start zender-worker:*
;
; Monitoring:
; - Check status: sudo supervisorctl status
; - View logs: tail -f /var/www/zender/storage/logs/worker.log
; ============================================

[program:zender-worker]
; Name pattern for multiple processes
process_name=%(program_name)s_%(process_num)02d

; Command to run the queue worker
; Optimizations:
;   --sleep=3      : Wait 3 seconds when no jobs (reduces CPU)
;   --tries=3      : Retry failed jobs 3 times
;   --max-jobs=1000: Restart after 1000 jobs (prevents memory leaks)
;   --max-time=3600: Restart every hour (prevents zombie processes)
;   --memory=128   : Restart if exceeds 128MB (memory safety)
command=php /var/www/zender/artisan queue:work redis --sleep=3 --tries=3 --max-jobs=1000 --max-time=3600 --memory=128

; Working directory
directory=/var/www/zender

; Auto-management
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

; Run as web server user
user=www-data

; Number of worker processes
; IMPORTANT: Set this based on your CPU cores
; For 2 vCPU: use 2
; For 4 vCPU: use 3-4
numprocs=2

; Logging
redirect_stderr=true
stdout_logfile=/var/www/zender/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; Grace period for jobs to finish before killing
stopwaitsecs=3600
